# Docker-Full-Image/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
user=root
loglevel=info
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:producer]
command=python /app/producer/coinbase_producer.py
autostart=true
autorestart=true
startretries=3
startsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/app"
user=root
priority=100
stopsignal=TERM
stopwaitsecs=30
killasgroup=true
stopasgroup=true

[program:consumer]
command=python /app/consumer/postgres_consumer.py
autostart=true
autorestart=true
startretries=3
startsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/app"
user=root
priority=200
stopsignal=TERM
stopwaitsecs=30
killasgroup=true
stopasgroup=true

[program:alert_consumer]
command=python /app/alerting/price_alert_consumer.py
autostart=true
autorestart=true
startretries=3
startsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/app"
user=root
priority=300
stopsignal=TERM
stopwaitsecs=30
killasgroup=true
stopasgroup=true

[program:throughput_monitor]
command=python /app/monitoring/throughput_monitor.py
autostart=true
autorestart=true
startretries=3
startsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/app"
user=root
priority=400
stopsignal=TERM
stopwaitsecs=30
killasgroup=true
stopasgroup=true

[group:crypto_pipeline]
programs=producer,consumer,alert_consumer,throughput_monitor
priority=999

[eventlistener:process_monitor]
command=python /app/monitoring/process_monitor.py
events=PROCESS_STATE_EXITED,PROCESS_STATE_STOPPED
autostart=true
autorestart=true
user=root


